# services/whisper-api.service
[Unit]
Description=Whisper Audio Transcription API Service
After=network.target

[Service]
Type=simple
User=www-data
Group=www-data
WorkingDirectory=/opt/whisperTest
# Conda 环境配置（部署脚本会自动更新这些路径）
Environment="PATH=/opt/conda/bin:/opt/conda/envs/whisper-api/bin"
# 如果使用 venv，取消注释下面这行并注释上面的 Environment
# Environment="PATH=/opt/whisperTest/venv/bin"
ExecStart=/opt/conda/envs/whisper-api/bin/gunicorn -c config/gunicorn_config.py src.api:api_app
# 如果使用 venv，取消注释下面这行并注释上面的 ExecStart
# ExecStart=/opt/whisperTest/venv/bin/gunicorn -c config/gunicorn_config.py src.api:api_app
Restart=always
RestartSec=10
StandardOutput=append:/opt/whisperTest/logs/api_service.log
StandardError=append:/opt/whisperTest/logs/api_service_error.log

# 资源限制
LimitNOFILE=65535
LimitNPROC=4096

# 安全设置
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target

